# com.filelay.icloud-file-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.filelay.icloud-file-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/Documents/VibeCoding/Filelay/icloud_file_sync.sh</string>
    <string>--source</string>
    <string>/Users/alex/Documents/your-local-file.txt</string>
    <string>--dest</string>
    <string>/Users/alex/Library/Mobile Documents/com~apple~CloudDocs/Sync/your-local-file.txt</string>
    <string>--interval</string>
    <string>1</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/filelay-icloud-sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/filelay-icloud-sync.err.log</string>
</dict>
</plist>
